# docs/plantcare.service
# =============================================================================
# Systemd Service Unit — Plant Care Monitoring Backend (FastAPI)
# =============================================================================
# Install on EC2:
#   sudo cp plantcare.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable plantcare
#   sudo systemctl start plantcare
# =============================================================================

[Unit]
Description=Plant Care Monitoring System — FastAPI Backend
After=network.target postgresql.service

[Service]
Type=simple
User=plantcare
Group=plantcare
WorkingDirectory=/opt/plantcare

# Run FastAPI with Uvicorn ASGI server
ExecStart=/opt/plantcare/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080

# Load environment variables from file
EnvironmentFile=/opt/plantcare/.env

# Restart policy
Restart=on-failure
RestartSec=10
SuccessExitStatus=143

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plantcare

[Install]
WantedBy=multi-user.target
